# system_manager/smarttv-local-system.service
[Unit]
Description=SmartTV Local System Management Server
After=network.target
StartLimitBurst=5
StartLimitIntervalSec=30

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/usr/local/share/smarttv
ExecStart=/usr/bin/python3 /usr/local/share/smarttv/local_system_server.py
Restart=always
RestartSec=5

# Environment variables
Environment=PYTHONUNBUFFERED=1

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=smarttv-local-system

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log

[Install]
WantedBy=multi-user.target